# evaluation/todo_setup/setup/todoapp.service
[Unit]
Description=Todo app, ACIT4640
After=network.target
Requires=mongod.service

[Service]
Environment=NODE_PORT=8080
WorkingDirectory=/home/todoapp/app
Type=simple
User=todoapp
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/node /home/todoapp/app/server.js
Restart=always

[Install]
WantedBy=multi-user.target